# systemd/carlottery-telegram-bot.service
[Unit]
Description=Car Lottery Telegram Bot
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/apps/CarLottery/backend

EnvironmentFile=/home/ubuntu/apps/CarLottery/.env

Environment="PATH=/home/ubuntu/apps/CarLottery/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/apps/CarLottery/backend"
Environment="DJANGO_SETTINGS_MODULE=bingo.settings"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/home/ubuntu/apps/CarLottery/venv/bin/python manage.py runbot

Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=carlottery-telegram-bot

NoNewPrivileges=true
PrivateTmp=true

TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
