# scripts/celery.service
[Unit]
Description=Celery worker daemon for Bingo game
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/apps/arif_bingo/backend
Environment="PATH=/home/ubuntu/apps/arif_bingo/backend/venv/bin"
Environment="PYTHONPATH=/home/ubuntu/apps/arif_bingo/backend"
Environment="DJANGO_SETTINGS_MODULE=bingo.settings"
ExecStart=/home/ubuntu/apps/arif_bingo/backend/venv/bin/celery -A bingo worker --concurrency=2 -Q gameplay,registration,default --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

